# systemd-monit/systemd/systemd-monit.service
[Service]
ExecStart=%h/.local/bin/systemd-monit
Type=oneshot

[Install]
WantedBy=multi-user.target
