# lego-wrapper/systemd/lego-wrapper@.timer
[Unit]
Description=Timer to run lego wrapper

[Timer]
OnCalendar=05:30
RandomizedDelaySec=30min

[Install]
WantedBy=timers.target
