# lego-wrapper/systemd/lego-wrapper@.service
[Unit]
Description=run lego wrapper

[Service]
Type=oneshot
ExecStart=/usr/bin/lego-wrapper %i

[Install]
WantedBy=default.target
