# update-notifier/systemd/update-notifier.service
[Service]
ExecStart=%h/.local/bin/update-notifier
Type=oneshot

[Install]
WantedBy=multi-user.target
