# update-notifier/systemd/update-notifier.timer
[Timer]
OnBootSec=60min
OnUnitActiveSec=1min

[Install]
WantedBy=timers.target
