# update-notifier/systemd/update-notifier.service
[Service]
ExecStart=update-notifier
Type=oneshot

[Install]
WantedBy=multi-user.target
